# packages/pglift/pglift-0.27.0.tar.gz/pglift-0.27.0/src/pglift/systemd/pglift-postgresql@.service
[Unit]
Description=PostgreSQL %i database server
After=network.target

[Service]
Type=forking
{executeas}

# Disable OOM kill on the postmaster
OOMScoreAdjust=-1000
Environment=PG_OOM_ADJUST_FILE=/proc/self/oom_score_adj
Environment=PG_OOM_ADJUST_VALUE=0

{environment}
ExecStart={python} -m pglift.postgres %i
ExecReload=/bin/kill -HUP $MAINPID

PIDFile={pid_directory}/postgresql-%i.pid

[Install]
WantedBy=multi-user.target
