# packages/pglift/pglift-0.27.0.tar.gz/pglift-0.27.0/src/pglift/systemd/pglift-backup@.timer
[Unit]
Description=Backup %i PostgreSQL database instance

[Timer]
OnCalendar={calendar}
Persistent=true

[Install]
WantedBy=timers.target
